# nginx/supervisor.conf
[program:airbook-uwsgi]
command = source /var/opt/airbook-server/environment.sh && /usr/local/bin/uwsgi -i /var/opt/airbook-server/nginx/uwsgi.ini
