# nginx/supervisor.conf
[program:airbook-uwsgi]
command = /usr/local/bin/uwsgi -i /var/opt/airbook-server/nginx/uwsgi.ini
